# pepaapi/configuracion/etc/supervisor/daemon-rs485area.conf
[program:daemon-rs485area]
directory=/var/www/html/pepaapi
command=php artisan command:Area54Daemon --debug
stdout_logfile=/var/www/html/pepaapi/storage/logs/daemons.log
redirect_stderr=true
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autostart=true
user=www-data
autorestart=true
